# infrastructure/systemd/cleancentive-reconcile.timer
[Unit]
Description=Fallback timer for Cleancentive production reconciliation

[Timer]
OnBootSec=10m
OnUnitActiveSec=60m
Persistent=true

[Install]
WantedBy=timers.target
